# nginx/dashboard.service
[Unit]
Description=POD Data Pipeline Dashboard API
After=network.target

[Service]
Type=simple
User=YOUR_SYSTEM_USER
WorkingDirectory=/opt/YOUR_BRAND_ID/dashboard
ExecStart=/usr/bin/python3 -m uvicorn api:app --host 127.0.0.1 --port 8090
Restart=always
RestartSec=5
EnvironmentFile=/opt/YOUR_BRAND_ID/.env

[Install]
WantedBy=multi-user.target
